# httpzipdir.service
[Service]
ExecStart=/path/to/httpzipdir --allow /URL/to/publish=/path/to/publish
DynamicUser=yes
